# dot_config/systemd/user/niri-session.target
[Unit]
Description=Niri session services
