# dot_config/systemd/user/noctalia-theme-sync.path
[Unit]
Description=Watch Noctalia colors and run theme sync

[Path]
PathChanged=%h/.config/noctalia/colors.json

[Install]
WantedBy=default.target
