# dot_config/systemd/user/noctalia-theme-sync-start.service
[Unit]
Description=Run Noctalia theme sync when the graphical session starts
After=niri.service
PartOf=graphical-session.target

[Service]
Type=oneshot
ExecStart=/usr/bin/bash -lc 'bash "$HOME/.config/noctalia/scripts/theme-sync.sh" --source session-start'

[Install]
WantedBy=graphical-session.target
