# dot_config/systemd/user/noctalia-theme-sync.service
[Unit]
Description=Run Noctalia theme sync after theme color changes

[Service]
Type=oneshot
ExecStart=/usr/bin/bash -lc 'bash "$HOME/.config/noctalia/scripts/theme-sync.sh" --source color-watcher'
